# ansible/roles/bitcoin-retry-endpoint/templates/bitcoin_retry_endpoint.rc.j2
#!/bin/sh
#
# PROVIDE: bitcoin_retry_endpoint
# REQUIRE: NETWORKING DAEMON
# KEYWORD: shutdown
#
# Managed by Ansible — do not edit manually
#
# Add the following to /etc/rc.conf to enable:
#   bitcoin_retry_endpoint_enable="YES"

. /etc/rc.subr

name="bitcoin_retry_endpoint"
rcvar="bitcoin_retry_endpoint_enable"
desc="Bitcoin Retry Endpoint"

load_rc_config "${name}"

: ${bitcoin_retry_endpoint_enable:="NO"}
: ${bitcoin_retry_endpoint_user:="{{ retry_user }}"}
: ${bitcoin_retry_endpoint_group:="{{ retry_group }}"}
: ${bitcoin_retry_endpoint_config:="/usr/local/etc/bitcoin-retry-endpoint.conf"}
: ${bitcoin_retry_endpoint_logfile:="/var/log/bitcoin_retry_endpoint.log"}
: ${bitcoin_retry_endpoint_pidfile:="/var/run/bitcoin_retry_endpoint.pid"}

pidfile="${bitcoin_retry_endpoint_pidfile}"
command="{{ retry_bin_dir }}/bitcoin-retry-endpoint"
command_args=""

start_precmd="bitcoin_retry_endpoint_prestart"
start_cmd="bitcoin_retry_endpoint_start"
stop_cmd="bitcoin_retry_endpoint_stop"
status_cmd="bitcoin_retry_endpoint_status"

bitcoin_retry_endpoint_prestart()
{
    # Load environment from config file
    if [ -f "${bitcoin_retry_endpoint_config}" ]; then
        . "${bitcoin_retry_endpoint_config}"
    fi
    export MULTICAST_IF LISTEN_PORT NACK_PORT
    export EGRESS_IFACE EGRESS_PORT SHARD_BITS MC_SCOPE MC_BASE_ADDR
    export CACHE_BACKEND REDIS_ADDR CACHE_TTL CACHE_MAX_KEYS
    export RL_IP_RATE RL_IP_BURST RL_SENDER_RATE RL_SENDER_BURST
    export RL_GLOBAL_RATE RL_GLOBAL_BURST
    export NUM_WORKERS METRICS_ADDR OTLP_ENDPOINT OTLP_INTERVAL
    export DRAIN_TIMEOUT INSTANCE_ID DEBUG
}

bitcoin_retry_endpoint_start()
{
    bitcoin_retry_endpoint_prestart
    /usr/sbin/daemon -f -p "${pidfile}" -u "${bitcoin_retry_endpoint_user}" \
        ${command} >> "${bitcoin_retry_endpoint_logfile}" 2>&1
    echo "bitcoin_retry_endpoint started."
}

bitcoin_retry_endpoint_stop()
{
    if [ -f "${pidfile}" ]; then
        kill -TERM $(cat "${pidfile}") && rm -f "${pidfile}"
        echo "bitcoin_retry_endpoint stopped."
    else
        echo "bitcoin_retry_endpoint is not running."
    fi
}

bitcoin_retry_endpoint_status()
{
    if [ -f "${pidfile}" ] && kill -0 $(cat "${pidfile}") 2>/dev/null; then
        echo "bitcoin_retry_endpoint is running (pid $(cat ${pidfile}))."
    else
        echo "bitcoin_retry_endpoint is not running."
        return 1
    fi
}

run_rc_command "$1"
